# bin/nxlucened
#!/bin/sh
#
# NXLucene:  Start/stop/restart/status service
#
# chkconfig:    2345 25 90
# description:  NXLucene
#
# $Id: $

RETVAL=0
PWD=`pwd`
NAME='nxlucened'

# Change me.
NXLUCENE_HOME='..'

export LD_LIBRARY_PATH=/usr/local/lib

PIDDIR=$NXLUCENE_HOME'/var'
LOGFILE=$NXLUCENE_HOME'/log/nxlucene.log'

function start() {
    echo "Starting NXLucene:"    
    export PYTHONPATH=$PYTHONPATH:$NXLUCENE_HOME'/src'
    exec $NXLUCENE_HOME'/bin/runnxlucene' &
}

function stop() {

    echo "Stoping NXLucene:"

    ps ax |grep twistd |grep $NAME &>/dev/null
    if [ "$?" -eq 0 ]; then
	kill -9 `cat $PIDDIR/$NAME.pid`
	rm -f $PIDDIR/$NAME.pid
    else
	echo
	echo -n "($NAME is not running)"
	rm -f $PIDDIR/$NAME.pid
    fi
}

function status() {
    if [ -s $PIDDIR/$NAME.pid ] ; then
	echo "NXLucene running..."
    else
	echo "NXLucene not running..."
    fi
}

function logtail() {
    tail -f $LOGFILE
}

# See how we were called.
case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
	status
        ;;
    restart)
        stop
        start
        ;;
    logtail)
	logtail
	;;
    *)
        echo $"Usage: $0 {start|stop|restart|status|logtail}"
        RETVAL=1
esac

exit $RETVAL
